# live-build/ubuntu-server-live/includes.binary/boot/overlay.dir/lib/systemd/system/getty@.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --noclear --login-pause --autologin ubuntu %I $TERM
