# live-build/ubuntu-server/includes.binary/overlay/lib/systemd/system/media-rack.lower.mount
[Mount]
What=/cdrom/casper/maas/maas-rack.squashfs
Where=/media/rack.lower
Type=squashfs
Options=ro
