# live-build/ubuntu-server/includes.chroot.ubuntu-server-minimal.ubuntu-server.installer/usr/lib/systemd/system/media-filesystem.mount
[Mount]
What=/cdrom/casper/ubuntu-server-minimal.ubuntu-server.squashfs
Where=/media/filesystem
Type=squashfs
Options=ro
