# live-build/ubuntu-server/includes.chroot.ubuntu-server-minimal.ubuntu-server.installer/usr/lib/systemd/system/media-filesystem.mount
[Unit]
Requires=media-full.mount media-minimal.mount
After=media-full.mount media-minimal.mount

[Mount]
What=/cdrom/casper/ubuntu-server-minimal.ubuntu-server.squashfs
Where=/media/filesystem
Type=overlay
Options=ro,lowerdir=/media/full:/media/minimal
