# live-build/ubuntu-server/includes.chroot.ubuntu-server-minimal.ubuntu-server.installer/usr/lib/systemd/system/media-minimal.mount
[Mount]
What=/cdrom/casper/ubuntu-server-minimal.squashfs
Where=/media/minimal
Type=squashfs
Options=ro
