# live-build/ubuntu-server/includes.binary/overlay/usr/lib/systemd/system/media-region.lower.mount
[Mount]
What=/cdrom/casper/maas/maas-region.squashfs
Where=/media/region.lower
Type=squashfs
Options=ro
