# live-build/ubuntu-server/includes.binary/overlay/usr/lib/systemd/system/subiquity_autoinstall.mount
[Mount]
What=/dev/disk/by-label/autoinstall
Where=/autoinstall
Type=ext4
