# live-build/ubuntu-server/includes.binary/overlay/usr/lib/systemd/system/serial-getty@.service.d/subiquity-serial.conf
[Unit]
Description=Subiquity, the installer for Ubuntu Server %I
After=snapd.seeded.service
StartLimitInterval=0

[Service]
Environment=SNAP_REEXEC=0
ExecStart=
ExecStart=/usr/bin/snap run subiquity.subiquity-service %I
