# live-build/ubuntu-server/includes.chroot.ubuntu-server-minimal.ubuntu-server.installer/usr/lib/systemd/system/serial-getty@.service.d/subiquity-serial.conf
[Unit]
Description=Subiquity, the installer for Ubuntu Server %I
StartLimitInterval=0

[Service]
Environment=SNAP_REEXEC=0
ExecStart=
ExecStart=/usr/bin/snap run subiquity.subiquity-service %I
