# live-build/ubuntu-server/includes.binary/overlay/usr/lib/systemd/system/snap.subiquity.subiquity-service.service.d/subiquity.conf
[Unit]
IgnoreOnIsolate=yes
After=systemd-user-sessions.service plymouth-quit-wait.service
ConditionPathExists=!/run/subiquity/complete

[Service]
Environment=SNAP_REEXEC=0
UtmpIdentifier=tty1
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
KillMode=process
IgnoreSIGPIPE=no
SendSIGHUP=yes
