# live-build/ubuntu-server/includes.binary/overlay/lib/systemd/system/getty@tty1.service.d/conditional.conf
[Unit]
ConditionPathExists=/run/subiquity/complete