# live-build/ubuntu-server/includes.binary/overlay/usr/lib/systemd/system/serial-getty@sclp_line0.service.d/subiquity-serial.conf
[Unit]
After=cloud-final.service snapd.seeded.service
[Service]
StandardOutput=tty
ExecStartPre=
ExecStartPre=/usr/bin/snap run subiquity --ssh
