# live-build/ubuntu-server/includes.binary/overlay/usr/lib/systemd/system/serial-getty@sclp_line0.service.d/subiquity-serial.conf
[Service]
ExecStart=
ExecStart=/usr/bin/snap run subiquity --ssh
